# services/input-mknod.service
[Unit]
Description=Input device node daemon (container isolation)
# Must run as system service (root) because mknod for device nodes
# requires CAP_MKNOD and rootful container to work

[Service]
Type=simple
ExecStart=/usr/local/bin/input-mknod-daemon
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
